# packages/backup-to-ssh/files/etc/systemd/system/backup-to-ssh.timer
[Unit]
Description=Backup to ssh host (timer)

[Timer]
OnCalendar=0:11:0
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target
